# systemd/gnostis.service
[Unit]
Description=Gnostis MCP HTTP server

[Service]
Type=simple
WorkingDirectory=%h/.gnostis
ExecStart=%h/.local/bin/gnostis run
Restart=unless-stopped
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
